# dbus.socket
[Unit]
Description=D-Bus Message Bus Socket

[Socket]
ListenStream=/run/user/%U/bus

[Install]
WantedBy=default.target
