# env.service
[Unit]
Description=User environment
Before=default.target

[Service]
Type=oneshot
ExecStart=-/usr/bin/systemctl --user set-environment XAUTHORITY=%h/.Xauthority GUILE_LOAD_COMPILED_PATH=%h/usr/lib/guile/2.0/site-ccache GUILE_LOAD_PATH=%h/usr/share/guile/site/2.0

[Install]
WantedBy=default.target
