# gpg-agent.service
[Unit]
Description=GPG Agent service
IgnoreOnIsolate=true
Before=default.target

[Service]
Type=forking
ExecStart=/usr/bin/gpg-agent --daemon

[Install]
WantedBy=default.target
