# graphical.target
[Unit]
Description=User Graphical Interface target

[Install]
WantedBy=default.target
