# tvguide.service
[Unit]
Description=TV guide service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tvguide-refresh
Nice=19

