# tvguide.timer
[Unit]
Description=TV guide timer

[Timer]
OnCalendar=Mon *-*-* 00:01
AccuracySec=10min

[Install]
WantedBy=default.target
