# wm.target
[Unit]
Description=WM (DE) target

[Install]
WantedBy=graphical.target
