# xsettings.target
[Unit]
Description=X defaults and settings target

[Install]
WantedBy=graphical.target
